# debian/postinst
#!/bin/sh

set -e

if [ -d /run/systemd/system ]; then
  systemctl enable sdr-server.service
	systemctl start sdr-server.service
fi